# sail-5.0.0/debian/sail/opt/sark/service/srk-ua-responder/run
#!/bin/sh

exec 2>&1
exec /opt/sark/scripts/responder.pl --nofork
