# sail-5.1.0/opt/sark/service/srk-ua-siplog/run
#!/bin/sh

exec 2>&1
exec /usr/sbin/tcpdump -U -i eth0 port 5060 -s 0 -w  /var/log/sip.log
