# sail-4.0.1/opt/sark/service/sark/run
#!/bin/sh
exec 2>&1
logger Starting Asterisk
exec nice -n 0 /usr/sbin/asterisk -f -U asterisk -G asterisk > /dev/null 2>&1
