# sail-6.0.1/debian/sail/etc/cron.d/sark
# SARK timer
* * * * * root php /opt/sark/generator/srktimer.php >/dev/null 2>&1 

# SARK queue mover in ARC
* * * * * root perl /opt/sark/scripts/srkcronqmove.pl >/dev/null 2>&1

# SARK spin off regression copies
0 4 * * * root sh /opt/sark/scripts/spin.sh >/dev/null 2>&1

# SARK age regressions
0 4 * * * root sh /opt/sark/scripts/age.sh >/dev/null 2>&1

#
# HD pruning disabled
#
# SARK HD pruning
#01 * * * * root perl /opt/sark/scripts/aelhdlon.pl prune >/dev/null 2>&1

# SARK VMAIL ageing
01 * * * * root perl /opt/sark/scripts/agevmail.sh >/dev/null 2>&1
#sark disk space pruner (ppc and arm)
*/10 * * * * root perl /opt/sark/scripts/diskabs.pl >/dev/null 2>&1

#
# Blacklist and whitelist loaders removed
# 
#sark blacklist loader for the firewall - requires ipset 
#*/30 * * * * root sh /opt/sark/scripts/getvoipbl.sh >/dev/null 2>&1
#sark fqdn whitelist loader for the firewall - requires ipset 
#* * * * * root /bin/sh /opt/sark/scripts/getfqdnip.sh 
#*/30 * * * * root sh /opt/sark/scripts/getfqdndrop.sh >/dev/null 2>&1

#run an internet speed test once every day
0 0 * * * root /usr/bin/speedtest-cli --simple  > /opt/sark/cache/speedtest 

#fetch manufacturer MAC list once a week
0 0 * * 0 root sh /opt/sark/scripts/getmaclist.sh

#refresh phone images  weekly
0 0 * * 0 root sh /opt/sark/scripts/getimages.sh

#recording offload every 10 minutes
*/10 * * * * root sh /opt/sark/scripts/srkoffload_simple.sh >/dev/null 2>&1
