# sail-7.0.0/opt/sark/service/srk-ua-siplog/run
#!/bin/sh

exec 2>&1
exec /usr/sbin/tcpdump -nqt -s 0 -A port 5060 -w  /var/log/siplog
